# infra/systemd/trading-bot.service
[Unit]
Description=Trading Bot Live Profile
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=gr8monk3ys
WorkingDirectory=/home/gr8monk3ys/trading-bot
EnvironmentFile=/home/gr8monk3ys/trading-bot/.env
ExecStart=/home/gr8monk3ys/trading-bot/.venv/bin/python scripts/run_low_resource_profile.py --asset-class crypto --risk-profile aggressive
Restart=always
RestartSec=10
StandardOutput=append:/home/gr8monk3ys/trading-bot/logs/systemd_trading_bot.log
StandardError=append:/home/gr8monk3ys/trading-bot/logs/systemd_trading_bot.log

[Install]
WantedBy=multi-user.target
